# contrib/pw-merger.service
[Unit]
Description=pw-merger — persistent PipeWire audio sink merger
Documentation=https://github.com/you/pw-merger
# Wait for the PipeWire socket to be available.
After=pipewire.service pipewire-pulse.service
Wants=pipewire.service

[Service]
Type=simple

# ── EDIT THESE TWO LINES ───────────────────────────────────────────────────
Environment=DEVICE_A=alsa_output.pci-0000_00_1f.3.analog-stereo
Environment=DEVICE_B=alsa_output.hdmi-stereo
# ──────────────────────────────────────────────────────────────────────────

Environment=SINK_NAME=Merged Output
Environment=RUST_LOG=pw_merger=info

ExecStart=%h/.local/bin/pw-merger \
    --sink-name  "${SINK_NAME}" \
    --device-a   "${DEVICE_A}" \
    --device-b   "${DEVICE_B}"

# Restart if the daemon crashes or PipeWire restarts.
Restart=on-failure
RestartSec=2s

# Give PipeWire 3 s to come up before we try to connect.
ExecStartPre=/bin/sleep 1

[Install]
WantedBy=default.target
